# deploy/store.service
[Unit]
Description=Stride Store Django application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/Store
EnvironmentFile=/var/www/Store/.env
ExecStart=/var/www/Store/.venv/bin/gunicorn store.wsgi:application --bind 127.0.0.1:8000 --workers 1 --timeout 120 --access-logfile - --error-logfile -
Restart=on-failure
RestartSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
